# supervisor.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:telegram_bot]
command=python main.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/telegram_bot.err.log
stdout_logfile=/var/log/supervisor/telegram_bot.out.log
environment=ENVIRONMENT=production

[program:admin_panel]
command=python run_admin.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/admin_panel.err.log
stdout_logfile=/var/log/supervisor/admin_panel.out.log
environment=ENVIRONMENT=production

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
